# systemd/telegram-files-sync.service
[Unit]
Description=Telegram Files Auto Sync Service
Documentation=https://rclone.org/
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/telegram-files-sync.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 安全设置
User=root
Group=root

# 资源限制
Nice=10
IOSchedulingClass=2
IOSchedulingPriority=7

# 环境变量
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# 超时设置
TimeoutStartSec=30
TimeoutStopSec=30

# 重启策略
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
